# tornado-git.conf
description "Tornado Git Server"
author "Connor Gorman"
  
start on (local-filesystems and net-device-up)
stop on runlevel [06]
  
respawn
exec /home/ubuntu/Webhook/tornado_git.py